# miniforge/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/chuan/Desktop/xiangmu/客户端/miniforge/bin/conda";
setenv _CONDA_EXE "/Users/chuan/Desktop/xiangmu/客户端/miniforge/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/Users/chuan/Desktop/xiangmu/客户端/miniforge/bin/python";
setenv _CONDA_ROOT "/Users/chuan/Desktop/xiangmu/客户端/miniforge";

source "/Users/chuan/Desktop/xiangmu/客户端/miniforge/lib/python3.12/site-packages/conda/shell/etc/profile.d/conda.csh"